# cron/crontab
# mosdns rule sync - daily at 00:00
0 0 * * * /etc/mosdns/cron/sync-rules.sh >> /etc/mosdns/tmp/sync-rules.log 2>&1

# mosdns candidate domain sync - hourly at minute 5, stats only, no reload
5 * * * * /etc/mosdns/cron/sync-candidates.sh >> /etc/mosdns/tmp/sync-candidates.log 2>&1

# mosdns candidate promotion - daily at 00:20
20 0 * * * /etc/mosdns/cron/promote-candidates.sh >> /etc/mosdns/tmp/promote-candidates.log 2>&1

# mosdns deferred reload/restart - daily at 00:30, at most once
30 0 * * * /etc/mosdns/cron/reload-if-needed.sh >> /etc/mosdns/tmp/reload.log 2>&1
